# server/deploy/ytdlp-api.service
[Unit]
Description=yt-dlp API Server
After=network.target

[Service]
Type=simple
User=ytdlp
Group=ytdlp
WorkingDirectory=/opt/ytdlp-api
ExecStart=/opt/ytdlp-api/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
EnvironmentFile=-/opt/ytdlp-api/.env
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/ytdlp-api

[Install]
WantedBy=multi-user.target
